# pidash.service
[Unit]
Description=PiDash Service

[Service]
ExecStart=bash -c '~/pidash/start.sh'

[Install]
WantedBy=multi-user.target